# endpoints/endpoints.service
# /etc/systemd/system/piper-endpoints.service
[Unit]
Description=Piper Endpoints API
After=network.target

[Service]
Type=simple
User=server
WorkingDirectory=/home/piper/endpoints
Environment="PATH=/home/piper/.venv/bin"
ExecStart=/home/piper/.venv/bin/python app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target